# docker/cron.d/scheduler
# Run a script every minute
* * * * * php /var/www/main/artisan schedule:run >> /dev/null 2>&1
* * * * * php /var/www/admin/artisan schedule:run >> /dev/null 2>&1
